# contrib/tbio-cron.timer
[Unit]
Description=Sync tbio

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
